# tests/scenarios/credential_stuffing.json
{
  "scenario_id": "TEST-CREDSTURF-001",
  "scenario_name": "Credential Stuffing Escalating to Lateral Movement",
  "description": "Baseline scenario matching the built-in SentinelOps demo. Simulates a credential stuffing attack against the login endpoint that escalates to lateral movement and privilege escalation.",
  "alert": {
    "alert_id": "ALERT-2025-0847",
    "alert_name": "Credential Stuffing Detection — Multiple Failed Logins",
    "alert_severity": "high",
    "alert_raw": {
      "source_ip": "40.80.148.42",
      "index": "sentinalops_os",
      "event_type": "credential_stuffing",
      "failed_count": 47,
      "success_count": 3,
      "targeted_accounts": ["svc_admin", "j.smith"],
      "compromised_hosts": ["web-prod-01", "app-srv-02", "we8105desk", "we9041srv"],
      "attack_stages": [
        "credential_stuffing",
        "initial_access",
        "lateral_movement",
        "privilege_escalation",
        "data_staging",
        "c2_communication",
        "persistence"
      ],
      "c2_domain": "hildegardsfarm.com",
      "data_staged_mb": 842
    }
  },
  "expected_blast_radius": {
    "compromised_entities": [
      {"type": "host", "name": "web-prod-01", "risk_score": 0.95},
      {"type": "host", "name": "app-srv-02", "risk_score": 0.88},
      {"type": "host", "name": "db-srv-01", "risk_score": 0.92},
      {"type": "user", "name": "svc_admin", "risk_score": 0.98},
      {"type": "user", "name": "j.smith", "risk_score": 0.75},
      {"type": "service", "name": "postgres-primary", "risk_score": 0.85},
      {"type": "data_store", "name": "customer_database", "risk_score": 0.90}
    ],
    "potentially_affected_entities": [
      {"type": "host", "name": "dc-01", "risk_score": 0.65},
      {"type": "service", "name": "nginx-frontend", "risk_score": 0.70},
      {"type": "service", "name": "api-backend", "risk_score": 0.72},
      {"type": "service", "name": "redis-cache", "risk_score": 0.55}
    ],
    "total_entity_count": 11,
    "compromised_count": 7,
    "business_impact": "HIGH — Service account compromised, lateral movement across 3 hosts, 842MB data staged for exfiltration, C2 channel active"
  },
  "expected_mitre_techniques": [
    "T1110.004 — Brute Force: Credential Stuffing",
    "T1078 — Valid Accounts",
    "T1021.001 — Remote Services: RDP",
    "T1570 — Lateral Tool Transfer",
    "T1569.002 — System Services: Service Execution",
    "T1098 — Account Manipulation",
    "T1071.001 — Application Layer Protocol: Web",
    "T1074.001 — Data Staged: Local Data Staging",
    "T1053.005 — Scheduled Task/Job"
  ],
  "expected_remediation_actions": [
    {
      "priority": 1,
      "title": "Block attacker IP 40.80.148.42 at perimeter firewall",
      "requires_approval": true
    },
    {
      "priority": 1,
      "title": "Force password reset for svc_admin and j.smith",
      "requires_approval": true
    },
    {
      "priority": 2,
      "title": "Isolate compromised hosts from production network",
      "requires_approval": true
    },
    {
      "priority": 2,
      "title": "Block C2 domain hildegardsfarm.com at DNS and proxy",
      "requires_approval": true
    },
    {
      "priority": 3,
      "title": "Enable MFA for all service accounts",
      "requires_approval": false
    },
    {
      "priority": 3,
      "title": "Remove scheduled task persistence mechanism",
      "requires_approval": true
    }
  ]
}
